# podsec-inotify/cron.hourly/podsec-inotify-check-policy
#!/bin/sh

/usr/lib/nagios/plugins/podsec-inotify-check-policy -vvv -c 100
