# podsec-inotify/crontabs/podsec-inotify-check-vuln
*/10    *       *       *       *       podsec-inotify-check-vuln
