# podsec-inotify/crontabs/podsec-inotify-check-policy
*/10    *       *       *       *       TMPF="/tmp/check-policy.$$";if /usr/lib/nagios/plugins/podsec-inotify-check-policy  -vvv -c 100 >$TMPF 2>&1; then :; else subject=$(head -1 $TMPF);  mail -s "$subject" root <$TMPF; rm -f $TMPF; fi
